# files/etc/uci-defaults/99-update-root-password
#!/bin/sh
# OpenWrt uci-defaults: 设置默认Web登录密码（root密码）
# 说明：
# 1) LuCI 的登录密码就是系统 root 的密码，因此只需改 /etc/shadow
# 2) 本脚本幂等；首启执行完会被系统自动删除

# ===【必改】把下行替换为你通过 "openssl passwd -6" 生成的哈希 ===
ROOT_HASH='$6$7wG/GFb3ktmCCNrt$PkSf3zltFde2/nva4PUJyRG6WX9rWu94P9tSp64Sr1bwbUoBNTZS/BiUpsS0cPzQqj0tWj4j6DzxM8mDeSdjB0'
# =================================================================

set -eu

update_root_password() {
    SHADOW='/etc/shadow'
    [ -f "$SHADOW" ] || exit 0
    # 仅当 root 的哈希不同才替换，避免重复写入
    if ! grep -q "^root:${ROOT_HASH}:" "$SHADOW" 2>/dev/null; then
        sed -i "s|^root:[^:]*:|root:${ROOT_HASH}:|" "$SHADOW"
        chmod 600 "$SHADOW" 2>/dev/null || true
    fi
}

update_root_password

exit 0
